# npm/.config/profile.d/50-npm.sh
NPM_CONFIG_USERCONFIG="$HOME/.config/npm/npmrc"
export NPM_CONFIG_USERCONFIG
PATH="$HOME/.local/npm/bin:$PATH"
export PATH
