# dpi/.config/autostart/dpi.desktop
[Desktop Entry]
Name=DPI
Comment=fuckfuckfuckfuckfugh
Exec=/usr/bin/auto-xrandr
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
